# flaskapp.service
[Unit]
Description=Gunicorn instance for a flask app
After=network.target
[Service]
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/LOG8415E
ExecStart=/home/ubuntu/LOG8415E/venv/bin/gunicorn -b localhost:8000 app:app
Restart=always
[Install]
WantedBy=multi-user.target